# mira.service
[Unit]
Description=Mira AI Copilot
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=piichan
WorkingDirectory=/home/piichan/mira
Environment=PYTHONUNBUFFERED=1
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/piichan/.Xauthority
ExecStart=/home/piichan/mira/venv/bin/python -m src.main
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical.target
